# prod.conf
# /etc/supervisor/conf.d

[inet_http_server]
port=127.0.0.1:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://127.0.0.1:9001

[supervisord]

[program:spot]
environment = PYTHONUNBUFFERED=1
user=sasthrayaan
directory=/home/sasthrayaan/sasthrayaan-spot
command=/home/sasthrayaan/sasthrayaan-spot/venv/bin/gunicorn --bind 0.0.0.0:8000 --workers 4 --log-level info config.wsgi
autostart=true
autorestart=true
stderr_logfile=/home/sasthrayaan/sasthrayaan-spot/log/gunicorn.err.log
stdout_logfile=/home/sasthrayaan/sasthrayaan-spot/gunicorn.out.log
